# libexec/rc/rc.d/vm_ddtrace_consumer
#!/bin/sh
#
#
# $FreeBSD$
#
# Start up for the Virtual DDTrace script reader daemon
#

# PROVIDE: vm_ddtrace_consumer
# REQURIE: syslogd
# BEFORE: DAEMON
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="vm_ddtrace_consumer"
desc="VM DDTrace consumer daemon"
# start_cmd="vm_ddtrace_consumer_start"
# stop_cmd=":"
command="usr/sbin/${name}"
# rcvar="vm_ddtrace_consumer_enable"

##vm_ddtrace_consumer_start()
# {
#    echo "This started."
#   ${vm_ddtrace_consumer_command}
#}

load_rc_config $name
run_rc_command "$1"